# deploy/crontab/privacyidea
# Run the privacyidea-cron runner every 5 minutes
*/5 * * * *	privacyidea	/opt/privacyidea/bin/privacyidea-cron run_scheduled -c
