# deploy/crontab/privacyidea
# Run the privacyidea-cron runner every 5 minutes
*/5 * * * *	privacyidea	/opt/privacyidea/bin/privacyidea-cron run_scheduled -c

# Audit rotation example: once a week at 00:00 on Sunday
#0 0 * * 0	privacyidea	/opt/privacyidea/bin/pi-manage audit rotate_audit --age 360

# Challenge cleanup: every day at 01:00
0 1 * * *    privacyidea     /opt/privacyidea/bin/pi-manage config challenge cleanup > /dev/null 2>&1
